# roles/reconciler/templates/substrate-reconcile-failure.service.j2
# Managed by substrate (roles/reconciler). Do not edit by hand.
[Unit]
Description=Substrate reconciler failure handler
Documentation={{ substrate_repo_url }}

[Service]
Type=oneshot
# substrate-reconcile.service triggers this via OnFailure= when a converge fails.
# systemd activates it on failure — it must NOT be enabled or started directly.
#
# The alert script reads a node-local webhook URL and POSTs a failure report; it
# skips loudly (and still exits 0) when no webhook is configured, so a node that
# opts out of alerting does not turn the failure handler itself into a failure.
# The leading '-' tells systemd to tolerate a non-zero exit anyway, so even a
# misbehaving alerter can never leave this unit failed and mask the original
# converge failure.
ExecStart=-{{ reconciler_alert_script }}
