# roles/reconciler/templates/substrate-reconcile.service.j2
# Managed by substrate (roles/reconciler). Do not edit by hand.
[Unit]
Description=Substrate reconciler (ansible-pull from main)
Documentation={{ substrate_repo_url }}
After=network-online.target
Wants=network-online.target
# Activated by systemd only when this unit fails; wires in alerting later.
OnFailure=substrate-reconcile-failure.service

[Service]
Type=oneshot
# --purge removes the checkout's stale files; --only-if-changed is intentionally
# NOT set so the node re-converges every run, repairing manual drift even when
# the spec on main has not changed.
ExecStart={{ reconciler_ansible_pull_bin }} \
    --url {{ substrate_repo_url }} \
    --checkout {{ substrate_branch }} \
    --directory {{ reconciler_workdir }} \
    --purge \
    --inventory localhost, \
    {{ substrate_playbook }}
# A failed converge must surface but must not wedge the node; the timer retries.
# Type=oneshot defaults to TimeoutStartSec=infinity, and the timer's
# OnUnitInactiveSec never fires while this unit is still activating — so one hung
# ansible-pull (e.g. a stalled git fetch) would permanently stop reconciliation.
# 10min exceeds any normal converge yet stays under the 15min timer interval.
TimeoutStartSec={{ reconciler_timeout_start_sec }}
# Record the outcome for fleet observability. $SERVICE_RESULT/$EXIT_STATUS are
# only set in ExecStopPost. Write atomically (tmp + mv) so readers never see a
# torn file. Literal percent signs are doubled (%%) because systemd expands %
# specifiers in unit files; Jinja passes % through untouched.
ExecStopPost=/bin/sh -c 'printf "last_run: %%s\nresult: %%s\nexit_status: %%s\n" "$(date -u +%%Y-%%m-%%dT%%H:%%M:%%SZ)" "$SERVICE_RESULT" "${EXIT_STATUS:-}" > {{ reconciler_status_file }}.tmp && mv {{ reconciler_status_file }}.tmp {{ reconciler_status_file }}'
Nice=10
IOSchedulingClass=best-effort
